# svrmg.conf
[program:svrmg]
command=/bin/bash -c "./start.sh"           ; the program (relative uses PATH, can take args)
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
directory=/project   ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
startretries=0
